# infrastructure/systemd/hosting-user-service.service
[Unit]
Description=Hosting Panel User Management Service (Go)
Documentation=https://github.com/xerudro/next-js-panel
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
Type=simple
User=hosting
Group=hosting
WorkingDirectory=/opt/hosting-panel/services/user-service

# Environment variables
Environment="GO_ENV=production"
EnvironmentFile=/opt/hosting-panel/services/user-service/.env

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hosting-panel/services/user-service/logs

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

# Start the service
ExecStart=/opt/hosting-panel/services/user-service/user-service
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=60
StartLimitBurst=3

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hosting-user-service

[Install]
WantedBy=multi-user.target
